# configuration_files/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --autologin iothermostat --noclear %I $TERM
Type=simple
